# systemd_user_service/templates/systemd-user.service.j2
[Unit]
Description={{ user_service.description }}

[Service]
ExecStart={{ user_service.exec }}

[Install]
WantedBy=multi-user.target
